# deploy/nimbo.service
[Unit]
Description=Nimbo — server management UI
Documentation=https://github.com/your-org/nimbo
# Start after the network is up. We *Want* Docker but do NOT *Require* it:
# the console must keep running (and be able to restart Docker) even if the
# Docker daemon dies. This is the whole point of running natively, not in a
# container.
After=network-online.target
Wants=network-online.target docker.service
# Disable the start-rate limiter so Restart=always keeps recovering a crash loop.
# (Must live in [Unit], not [Service], or systemd silently ignores it.)
StartLimitIntervalSec=0

[Service]
Type=simple
# Runs as the dedicated `nimbo` account (created by install.sh) with passwordless
# sudo. Privileged commands run via sudo (NIMBO_SUDO=1 in the env file), so the
# service process itself is not root.
User=nimbo
Group=nimbo
WorkingDirectory=/opt/nimbo
EnvironmentFile=/etc/nimbo/nimbo.env
ExecStart=/usr/bin/node /opt/nimbo/server.js

# Resilience: always come back, even after crashes or OOM.
Restart=always
RestartSec=2

# Logs go to the journal: journalctl -u nimbo -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nimbo

# Light hardening (kept permissive because the app needs host access).
NoNewPrivileges=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
